# sysinfo/files/etc/uci-defaults/99-sysinfo
#!/bin/sh

if grep -q "banner.failsafe" /etc/profile; then
	sed -i 's!/etc/banner.failsafe!/etc/banner.failsafe || /sbin/sysinfo.sh!' $TARGET_DIR/etc/profile
fi

exit 0
