# plugin-gargoyle-smsbox/files/etc/uci-defaults/plugin-gargoyle-smsbox
#!/bin/sh

uci set gargoyle.display.system_smsbox='SMS'
uci set gargoyle.scripts.system_smsbox='smsbox.sh'
uci set gargoyle.system.smsbox='351'
uci commit gargoyle
