# easyconfig/files/etc/uci-defaults/zz-easyconfig.sh
#!/bin/sh
[ ! -e /www/index.html ] && ln -s /www/easyconfig.html /www/index.html
exit 0
