# easyconfig/files/etc/hotplug.d/dhcp/99-easyconfig
#!/bin/sh
case "$ACTION" in
	add | update | remove)
		/usr/bin/easyconfig_wlanlogs.sh
		;;
esac
