# plugin-gargoyle-ftp/files/etc/uci-defaults/plugin-gargoyle-ftp
#!/bin/sh

uci set gargoyle.display.system_ftp='Serwer FTP'
uci set gargoyle.scripts.system_ftp='ftp.sh'
uci set gargoyle.system.ftp='351'
uci commit
