# easyconfig/files/etc/hotplug.d/iface/99-easyconfig
#!/bin/sh

[ $ACTION = "ifup" ] || exit 0

[ $INTERFACE = "wan" ] && awk '{printf "%d\n", $1}' /proc/uptime >> /tmp/easyconfig_wan_up

exit 0
